# root/etc/services.d/ortha/run
#!/usr/bin/with-contenv sh

echo "Starting orthanc..."
s6-setuidgid root /usr/local/sbin/Orthanc /etc/orthanc/orthanc.json
